# homelab-barman-service/image/crontab
* * * * * /usr/bin/barman cron --keep-descriptors
